# etc/root.crontab
### System ###
30  0  *   *   *    /bin/bash -l /usr/local/src/k8s/k8s/bin/cronjobs/system/system_report.sh
30  5  *   *   *    /bin/bash -l -c "/usr/local/src/k8s/k8s/bin/cronjobs/system/clean_backups.py"

### Nextcloud ###
# Full backup
30  4   *   *   6   /bin/bash -l /usr/local/src/k8s/k8s/bin/cronjobs/nextcloud/full_backup.sh
# Incremental backup
30  4   *   *   1-5 /bin/bash -l /usr/local/src/k8s/k8s/bin/cronjobs/nextcloud/incremental_backup.sh

### Palworld ###
30  10  *   *   *   mv /data/k8s/pv/palworld/backups/* /mnt/hdd-backup/backups/k8s/pv/palworld

### Growi ###
# Full backup
0   5   *   *   6   /bin/bash -l /usr/local/src/k8s/k8s/bin/cronjobs/growi/full_backup.sh
# Incremental backup
0   5   *   *   1-5 /bin/bash -l /usr/local/src/k8s/k8s/bin/cronjobs/growi/incremental_backup.sh
